# deploy/systemd/arbi-backend.service
[Unit]
Description=Gold MCX Arbitrage Backend (FastAPI + Dhan feed)
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
User=vs.bitcoding
WorkingDirectory=/opt/arbi/backend
EnvironmentFile=/opt/arbi/backend/.env
ExecStart=/opt/arbi/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
StandardOutput=append:/var/log/arbi/backend.log
StandardError=append:/var/log/arbi/backend.err.log

[Install]
WantedBy=multi-user.target
