# deploy/systemd/arbi-backend.service
[Unit]
Description=Gold MCX Arbitrage Backend (FastAPI + Dhan feed)
After=network.target

[Service]
Type=simple
User=vs.bitcoding
WorkingDirectory=/home/vs.bitcoding/gold-mcx-arbitrage/backend
EnvironmentFile=/home/vs.bitcoding/gold-mcx-arbitrage/backend/.env
ExecStart=/home/vs.bitcoding/gold-mcx-arbitrage/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
